# spectrum/bin/psm.monkey.service
[Unit]
Description=Pi Spectrum Monitor Monkey Service
After=multi-user.target
Wants=network-online.target
After=network-online.target 

[Service]
Type=idle
ExecStart=/usr/local/bin/psm-monkey
Restart=always
RestartSec=1
WorkingDirectory=~
RuntimeDirectory=psm
User=ses
Group=ses

[Install]
WantedBy=multi-user.target
