# spectrum/bin/psm.modprobe.service
[Unit]
Description=Create /dev/dsp* for audio devices

[Service]
Type=oneshot
ExecStart=/sbin/modprobe snd_pcm_oss

[Install]
WantedBy=multi-user.target
