# spectrum/bin/psm.email.service
[Unit]
Description=Send email to central server

[Service]
Type=oneshot
ExecStart=/home/ses/spectrum/piSpecMon/spectrum/bin/mail.sh

[Install]
WantedBy=multi-user.target
