# spectrum/bin/psm.target
[Unit]
Description=piSpecMon
After=basic.target

[Install]
WantedBy=basic.target
