# spectrum/bin/psm.monkey.service
[Unit]
Description=Pi Spectrum Monitor Monkey Service
After=multi-user.target
Wants=network-online.target psm.elasticsearch.service
After=network-online.target psm.elasticsearch.service

[Service]
Type=idle
ExecStart=/usr/local/bin/psm-monkey
Restart=always
RestartSec=1
WorkingDirectory=~
User=ses
Group=pi

[Install]
WantedBy=multi-user.target
