# spectrum/bin/psm.pico.service
[Unit]
Description=picofssd script
After=multi-user.target

[Service]
Type=idle
ExecStart=/usr/bin/python -u /usr/local/bin/picofssd.py
Restart=always
RestartSec=1
WorkingDirectory=~

User=root
Group=root

[Install]
WantedBy=multi-user.target
