# spectrum/bin/psm.hamlib-worker.service
[Unit]
Description=Pi Spectrum Monitor Worker Service
After=basic.target

[Service]
Type=idle
ExecStart=/usr/local/bin/psm-hamlib-worker
Restart=always
RestartSec=1
WorkingDirectory=~
RuntimeDirectory=psm-hamlib-worker
User=ses
Group=ses

[Install]
WantedBy=psm.target
