# spectrum/bin/psm.pico.service
[Unit]
Description=picofssd script

[Service]
Type=oneshot
ExecStart=/usr/bin/python -u /usr/local/bin/picofssd.py

[Install]
WantedBy=multi-user.target
