# spectrum/bin/psm.event.service
[Unit]
Description=Pi Spectrum Monitor Event Manager
After=basic.target

[Service]
Type=idle
ExecStart=/usr/local/bin/psm-event
Restart=always
RestartSec=1
WorkingDirectory=~
User=ses
Group=ses

[Install]
WantedBy=psm.target
