# spectrum/bin/psm.monkey.service
[Unit]
Description=Pi Spectrum Monitor Monkey Service
After=basic.target

[Service]
Type=idle
ExecStart=/usr/local/bin/psm-monkey
Restart=always
RestartSec=1
WorkingDirectory=~
RuntimeDirectory=psm-monkey
User=ses
Group=ses

[Install]
WantedBy=psm.target
