# spectrum/bin/psm.worker.service
[Unit]
Description=Pi Spectrum Monitor Worker Service
After=multi-user.target
Wants=network-online.target psm.elasticsearch.service
After=network-online.target psm.elasticsearch.service

[Service]
Type=idle
ExecStart=/usr/bin/python -u /home/ses/spectrum/piSpecMon/spectrum/worker.py
Restart=always
RestartSec=1
WorkingDirectory=/home/ses/spectrum/piSpecMon/spectrum

[Install]
WantedBy=multi-user.target
