# spectrum/bin/psm.worker.service
[Unit]
Description=Pi Spectrum Monitor Worker Service
After=multi-user.target
Wants=network-online.target
After=network-online.target

[Service]
Type=idle
ExecStart=/usr/local/bin/psm-worker
Restart=always
RestartSec=1
WorkingDirectory=~
RuntimeDirectory=psm
User=ses
Group=ses

[Install]
WantedBy=multi-user.target
