# spectrum/bin/psm.wav2mp3.service
[Unit]
Description=Pi Spectrum Monitor wav To mp3 Conversion Service
After=basic.target

[Service]
Type=idle
Environment=SES_CONFIG_PATH=/etc/psm.yml
ExecStart=/usr/bin/python -u /usr/local/bin/psm-wav2mp3
Restart=always
RestartSec=1
WorkingDirectory=~

User=ses
Group=ses

[Install]
WantedBy=psm.target
