# spectrum/bin/psm.rdevice.service
[Unit]
Description=Pi Spectrum Monitor RDevice Service
After=basic.target

[Service]
Type=idle
Environment=SES_CONFIG_PATH=/etc/psm.yml
ExecStart=/usr/local/bin/psm-rdevice
Restart=always
RestartSec=1
WorkingDirectory=~
RuntimeDirectory=psm-rdevice
User=ses
Group=ses

[Install]
WantedBy=psm.target
