# spectrum/bin/psm.3g.service
[Unit]
Description=3G modem driver
After=multi-user.target

[Service]
Type=idle
ExecStart=/home/ses/spectrum/piSpecMon/spectrum/bin/umtskeeper.sh
Restart=always
RestartSec=1
WorkingDirectory=/home/ses/spectrum/piSpecMon/spectrum

User=ses
Group=pi

[Install]
WantedBy=multi-user.target
