# spectrum/bin/psm.pico.service
[Unit]
Description=picofssd script
After=basic.target

[Service]
Type=idle
ExecStart=/usr/bin/python -u /usr/local/bin/picofssd.py
Restart=always
RestartSec=1
WorkingDirectory=~

User=ses
Group=ses

[Install]
WantedBy=psm.target
