# spectrum/bin/psm.audio.service
[Unit]
Description=Create /dev/dsp* for audio devices
After=sound.target

[Service]
Type=oneshot
ExecStart=/home/ses/piSpecMon/spectrum/bin/audio.sh

[Install]
WantedBy=psm.target
