# spectrum/bin/psm.email.service
[Unit]
Description=Send email to central server
After=network-online.target

[Service]
Type=oneshot
ExecStart=/home/ses/piSpecMon/spectrum/bin/mail.sh

[Install]
WantedBy=psm.target
