# spectrum/bin/psm.ams_worker.service
[Unit]
Description=Pi Spectrum Monitor AMS Sensor Worker Service
After=basic.target

[Service]
Type=idle
ExecStart=/usr/local/bin/psm-ams-worker
Restart=always
RestartSec=1
WorkingDirectory=~
RuntimeDirectory=psm-ams-worker
User=root
Group=root

[Install]
WantedBy=psm.target
