# spectrum/bin/psm.3g.service
[Unit]
Description=3G modem driver
After=network.target

[Service]
Type=idle
ExecStart=/home/ses/piSpecMon/spectrum/bin/sakis3g.sh connect
KillMode=process
Restart=always
RestartSec=60

[Install]
WantedBy=psm.target
